# services/transmission/run
#!/bin/bash
umask 0002
exec /sbin/setuser abc /usr/bin/transmission-daemon -g /config

